# roles/docker/templates/docker-prune.service.j2
[Unit]
Description=Docker system prune with volumes

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune --volumes -f
